# contrib/cliphist.service
[Unit]
Description=Wayland clipboard manager with support for multimedia 
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/wl-paste --watch cliphist store
Restart=on-failure

[Install]
WantedBy=graphical-session.target

# vim:ft=dosini
